# red/systemd-rm-tmp-files.service
[Unit]
Description=systemd-rm-tmp-files service
After=network-online.target

[Service]
ExecStart=/path

[Install]
WantedBy=multi-user.target